# launchd/com.alcatraz.claude-ipc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  claude-ipc broker — a KeepAlive LaunchAgent (a persistent daemon, NOT a
  scheduled cron, so no companion calendar event is required; observe it via
  `claude-ipc daemon status` + `claude-ipc tail`).

  Install:   cp launchd/com.alcatraz.claude-ipc.plist ~/Library/LaunchAgents/
             launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.alcatraz.claude-ipc.plist
  Uninstall: launchctl bootout gui/$(id -u)/com.alcatraz.claude-ipc
  (legacy `launchctl load/unload` is deprecated on current macOS — it fails with
   "Load failed: 5: Input/output error"; use bootstrap/bootout above.)

  Runs the COMPILED binary (`dist/claude-ipc serve`), not bun-on-source, so the
  broker and the CLI are the same artifact — `bun run build` updates both at once
  and there's no source-vs-dist drift. Rebuild before restarting after a change.
  Adjust the absolute path below if the repo moves.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.alcatraz.claude-ipc</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/alcatraz627/Code/Claude/claude-ipc/dist/claude-ipc</string>
    <string>serve</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/alcatraz627/.claude-ipc/logs/broker.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/alcatraz627/.claude-ipc/logs/broker.err.log</string>
</dict>
</plist>
